# ddns/sources/duckdns.service
[Unit]
Description=Run DuckDNS DDNS synchronization script

[Service]
ExecStart=/usr/local/bin/duckdns.sh
